# linux-vm-docker-dotnetcore-nodejs/scripts/crontab
# cleanup build jobs.
0 0 * * * /usr/bin/find /mnt/build/agent/_work -maxdepth 1 -mtime +30 -type d -regextype egrep -regex '.*/[0-9]{1,3}' -exec rm -rf {} +
